# new_project/learn/learn-platform.service
[Unit]
Description=Learning Platform Backend (LoongArch)
Documentation=https://github.com/spring-projects/spring-boot
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=forking
User=learn
Group=learn
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-loongarch64
Environment=SPRING_PROFILES_ACTIVE=loongarch
WorkingDirectory=/opt/learn-platform
ExecStart=/opt/learn-platform/scripts/start.sh
ExecStop=/opt/learn-platform/scripts/stop.sh
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/opt/learn-platform/app.pid
KillMode=process
Restart=on-failure
RestartSec=10s
TimeoutStartSec=60s
TimeoutStopSec=30s

# 资源限制（针对龙芯4核8GB配置）
LimitNOFILE=65535
LimitNPROC=32768
LimitMEMLOCK=infinity

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/learn-platform/logs
ReadWritePaths=/opt/learn-platform/app.pid

[Install]
WantedBy=multi-user.target
Alias=learn-platform.service 